# windows/profile.d/ansible
#ansible-webhook - AWX-like webhook for firstboot some server after installation finished
ansible-webhook
ansible-hook
