# windows/profile.d/completion.ps1
if (!(Test-Path HKCR:)) {
    $null = New-PSDrive -Name HKCR -PSProvider Registry -Root HKEY_CLASSES_ROOT
    $null = New-PSDrive -Name HKU -PSProvider Registry -Root HKEY_USERS
}

$scriptblock = {
    param($wordToComplete, $commandAst, $cursorPosition)
    dotnet complete --position $cursorPosition $commandAst.ToString() |
        ForEach-Object {
            [System.Management.Automation.CompletionResult]::new($_, $_, 'ParameterValue', $_)
        }
}
Get-Command -Name dotnet -ErrorAction SilentlyContinue -ErrorVariable dotnet_err -OutVariable dotnet_out > $null
if ($?) {
    Register-ArgumentCompleter -Native -CommandName dotnet -ScriptBlock $scriptblock
}