# windows/profile.d/exports.ps1
$Env:EDITOR = 'nvim.exe'
$Env:HOSTS = "${Env:WINDIR}\System32\drivers\etc\hosts"
